# debian/lttng-tools.lttng-sessiond.upstart
description "LTTng 2.0 central tracing registry session daemon"
author "Stéphane Graber <stgraber@ubuntu.com>"

start on local-filesystems
stop on runlevel [06]

respawn

script
	test -f /etc/default/lttng-sessiond && . /etc/default/lttng-sessiond
	exec /usr/bin/lttng-sessiond $DAEMON_ARGS
end script
