# debian/lttng-tools.lttng-sessiond.service
[Service]
ExecStart=/usr/bin/lttng-sessiond
