# debian/lttng-tools.lttng-sessiond.service
[Unit]
Description=LTTng session daemon
Documentation=http://lttng.org/docs/v2.7/#doc-lttng-sessiond

[Service]
ExecStart=/usr/bin/lttng-sessiond
